# tsdb/templates/default/startscript.sh.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          testone
# Required-Start:    $local_fs
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# X-Interactive:     false
# Short-Description: Example init script
# Description:       Start/stop an example script
### END INIT INFO
DESC="test script"
NAME=testone
#DAEMON=

do_start()
{
   echo "starting!";
   cd /home/opentsdb/build
 echo -ne '\n'| nohup ./tsdb tsd --port=4242 --staticroot=/home/opentsdb/build/staticroot --cachedir=/tmp/cache1 --zkquorum=<%= node[:hadoop][:zookeeper]%>:2181 --auto-metric &
}

do_stop()
{
   echo "stopping!"
   ps aux | grep 'opentsdb' | head -1 | awk '{print $2}' | xargs kill
}


case "$1" in
   start)
     do_start
     ;;
   stop)
     do_stop
     ;;
esac

exit 0